# cron/fetchavgblocktime.cron
0 0 * * * /usr/bin/php /var/www/bitcoinblockhalf.com/getavgblocktime.php