# cron/fetchprice.cron
*/15 * * * * /usr/bin/php /var/www/bitcoinblockhalf.com/getprice.php